# launchd/com.github.plaonn.hotspot-proxy-toggle.menu.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.github.plaonn.hotspot-proxy-toggle.menu</string>
  <key>ProgramArguments</key>
  <array>
    <string>__MENU_BIN__</string>
    <string>--command</string>
    <string>__INSTALL_BIN__</string>
    <string>--state</string>
    <string>__UI_STATE_PATH__</string>
    <string>--refresh</string>
    <string>__MENU_REFRESH_SECONDS__</string>
    <string>--title</string>
    <string>__MENU_BAR_TITLE__</string>
    <string>--locale</string>
    <string>__MENU_BAR_LOCALE__</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/hotspot-proxy-toggle-menu.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/hotspot-proxy-toggle-menu.stderr.log</string>
</dict>
</plist>
